# deploy/shuxue-worker.service
[Unit]
Description=shuxue.icu AI Task Worker
After=network.target postgresql.service shuxue.service
Requires=shuxue.service

[Service]
User=nginx
Group=nginx
WorkingDirectory=/var/www/shuxue/website/scripts
Environment="PATH=/var/www/shuxue/venv311/bin:/usr/bin"
EnvironmentFile=-/etc/shuxue/shuxue.env
ExecStart=/var/www/shuxue/venv311/bin/python worker.py
KillSignal=SIGTERM
TimeoutStopSec=180
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
